# files/system/etc/systemd/system/nix-daemon.socket.d/override.conf
[Unit]
After=nix.mount
BindsTo=nix.mount